# docker/app/config/supervisor/programs.conf
[program:cron]
command=crond -f -l 2
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=3
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stdout_logfile_maxbytes=50MB    ; Rotate stdout log when it reaches 50MB
stdout_logfile_backups=10       ; Keep 10 backup copies of stdout log
stderr_logfile=/var/log/supervisor/%(program_name)s.error.log
stderr_logfile_maxbytes=50MB    ; Rotate stderr log when it reaches 50MB
stderr_logfile_backups=10       ; Keep 10 backup copies of stderr log

[program:php-fpm]
command=php-fpm --nodaemonize
autorestart=true
autostart=true
startsecs=5
startretries=1
stopwaitsecs=3
killasgroup=true
stopasgroup=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stdout_logfile_maxbytes=50MB    ; Rotate stdout log when it reaches 50MB
stdout_logfile_backups=10       ; Keep 10 backup copies of stdout log
stderr_logfile=/var/log/supervisor/%(program_name)s.error.log
stderr_logfile_maxbytes=50MB    ; Rotate stderr log when it reaches 50MB
stderr_logfile_backups=10       ; Keep 10 backup copies of stderr log

; Symfony Scheduler worker: dispatches FetchRatesMessage every 5 min and
; CheckFeedIntegrityMessage hourly.
[program:scheduler-rates]
command=php -d xdebug.mode=off /app/bin/console messenger:consume scheduler_rates --time-limit=3600 --memory-limit=128M
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/%(program_name)s.error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10